# .crontab
# *    *    *    *    *
# ┬    ┬    ┬    ┬    ┬
# │    │    │    │    └─  Weekday  (0=Sun .. 6=Sat)
# │    │    │    └──────  Month    (1..12)
# │    │    └───────────  Day      (1..31)
# │    └────────────────  Hour     (0..23)
# └─────────────────────  Minute   (0..59)
#
0 10 * * 2,5 cd ~/dotfiles && ./scripts/autoupdate.sh
0 19 1,15 * * cd ~/dotfiles && ./scripts/map-dl.sh
0 9-21 * * * cd ~/dotfiles && ./scripts/n8n.sh >/dev/null 2>&1
